# core/sigilgate-sync.service
[Unit]
Description=SigilGate repository sync
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=<SERVICE_USER>
EnvironmentFile=<SCRIPTS_PATH>/.env
ExecStart=<SCRIPTS_PATH>/core/sync.sh
