# core/sigilgate-sync.timer
[Unit]
Description=SigilGate repository sync timer

[Timer]
OnCalendar=*:0/15
RandomizedDelaySec=60
Persistent=true

[Install]
WantedBy=timers.target
